# packaging/deb/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable nextnfs.service
systemctl start nextnfs.service || true
